# tmp/crontab
00 19 4 2 * usuario /ubicacion/del/script/consulta.sh
00 19 4 2 * usuario /ubicacion/del/script/consulta.sh
00 19 4 2 * usuario /ubicacion/del/script/consulta.sh
00 19 4 2 * usuario /ubicacion/del/script/consulta.sh
00 19 4 2 * usuario /ubicacion/del/script/consulta.sh
